# systemd-service-units-configuration/openbts.service
[Unit]
Description=OpenBTS

[Service]
Type=simple
ExecStart=/OpenBTS/OpenBTS
WorkingDirectory=/OpenBTS

[Install]
WantedBy=multi-user.target
