# systemd-service-units-configuration/sipauthserve.service
[Unit]
Description=OpenBTS SIP Authorization Service

[Service]
Type=simple
ExecStart=/OpenBTS/sipauthserve

[Install]
WantedBy=multi-user.target
